# scripts/service/nxctl-daemon.service
[Unit]
Description=NXCTL Daemon and API
After=docker.service

[Service]
Type=simple
User=root
# We use --with-api to run both monitoring and API in one service
ExecStart=/usr/local/bin/nxctl daemon --interval 10 --with-api
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
